# dist/linux/systemd/joyshockmapper@.service
[Unit]
Description=JoyShockMapper - Controller Remapper
Documentation=https://github.com/aapokaapo/JoyShockMapper
PartOf=graphical-session.target
After=graphical-session.target gnome-session.target

[Service]
Type=notify
# Wait up to 30 seconds for the AppIndicator extension's D-Bus interface to appear
ExecStartPre=/bin/bash -c 'for i in {1..30}; do busctl --user list | grep -q "org.kde.StatusNotifierWatcher" && exit 0; sleep 1; done; exit 0'
ExecStart=/usr/local/bin/JoyShockMapper --headless
ExecReload=/bin/kill -HUP $MAINPID

# Forward the user's display environment so the tray icon works correctly on both X11 and Wayland.
# AutoLoad (window-focus config switching) is disabled in headless mode.
PassEnvironment=DISPLAY WAYLAND_DISPLAY DBUS_SESSION_BUS_ADDRESS XDG_RUNTIME_DIR XDG_CONFIG_HOME HOME

# Restart on crash with exponential-style backoff to prevent restart loops
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=joyshockmapper

# Timeout & Watchdog
TimeoutStartSec=30
TimeoutStopSec=10
WatchdogSec=30

# Resource limits
MemoryMax=256M

[Install]
WantedBy=graphical-session.target
